# systemd/dark-matrix.service
[Unit]
Description=dark-matrix LED matrix daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/share/dark-matrix/node %h/.local/share/dark-matrix/daemon.mjs
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
